# life-management.service
[Unit]
Description=自律人生 - 个人生活管理系统
After=network.target

[Service]
Type=simple
User=life-user
Group=life-user
WorkingDirectory=/opt/life-management
ExecStart=/opt/life-management/自律人生-linux
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=life-management

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/life-management/data

[Install]
WantedBy=multi-user.target